# conf/supervisor/production.ini
[program:gunicorn]
command=/home/ubuntu/environments/Residency/bin/gunicorn wsgi:application -c /var/www/Residency/conf/gunicorn/current.py
directory=/var/www/Residency
user=nobody
group=nogroup
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/var/log/gunicorn.log

